# conf/systemd.service
[Unit]
Description=Lingua 90 language-learning server (static site + profile API)
After=network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
WorkingDirectory=__INSTALL_DIR__/

# Internal reverse-proxy port assigned by YunoHost
Environment=PORT=__PORT__
Environment=HOST=127.0.0.1
Environment=NODE_ENV=production
# "auto" => use the YunoHost-logged-in user (YNH_USER / Authorization header)
# as the learner profile when present; otherwise fall back to open profiles.
Environment=SSO_MODE=auto
# Per-learner JSON profiles live here (survives upgrades)
Environment=DATA_DIR=__DATA_DIR__

ExecStart=/usr/bin/node __INSTALL_DIR__/server/index.mjs

StandardOutput=append:/var/log/__APP__/__APP__.log
StandardError=inherit
Restart=on-failure
RestartSec=5

# Sandboxing — the service only needs to read its install dir and write its
# data dir + log.
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=__DATA_DIR__ /var/log/__APP__
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes
ProtectClock=yes
ProtectHostname=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
